# apache2.conf
[program:apache2]
command=/my_service apache2
